# vdsm_hooks/checkips/vdsm-checkips.service.in
# SPDX-FileCopyrightText: Red Hat, Inc.
# SPDX-License-Identifier: GPL-2.0-or-later

[Unit]
Description=Service that ping given ips each 60 seconds

[Service]
Type=simple
LimitCORE=infinity
ExecStart=@HOOKSDIR@/checkipsd
Restart=on-failure
RestartSec=10
User=@VDSMUSER@
Group=@VDSMGROUP@
TimeoutStopSec=@SERVICE_STOP_TIMEOUT@

[Install]
WantedBy=multi-user.target
