# static/usr/lib/systemd/system/mom-vdsm.service.in
# SPDX-FileCopyrightText: Red Hat, Inc.
# SPDX-License-Identifier: GPL-2.0-or-later

[Unit]
Description=MOM instance configured for VDSM purposes
Requires=vdsmd.service
After=vdsmd.service
Conflicts=momd.service ksmtuned.service

[Service]
Type=simple
LimitCORE=infinity
ExecStart=/usr/sbin/momd -c /etc/vdsm/mom.conf
Restart=always
RestartSec=5
User=@VDSMUSER@
Group=@VDSMGROUP@
TimeoutStopSec=@SERVICE_STOP_TIMEOUT@
StartLimitBurst=10
StartLimitInterval=120

[Install]
WantedBy=multi-user.target
