# static/usr/lib/systemd/system/supervdsmd.service.in
[Unit]
Description=Auxiliary vdsm service for running helper functions as root
Requires=libvirtd.service
Wants=lldpad.service
After=lldpad.service libvirtd.service

[Service]
Type=simple
LimitCORE=infinity
EnvironmentFile=-/etc/sysconfig/supervdsmd
ExecStart=@VDSMDIR@/daemonAdapter "@VDSMDIR@/supervdsmd" --sockfile "@VDSMRUNDIR@/svdsm.sock"
Restart=always
