# static/usr/lib/systemd/system/supervdsmd.service.in
# SPDX-FileCopyrightText: Red Hat, Inc.
# SPDX-License-Identifier: GPL-2.0-or-later

[Unit]
Description=Auxiliary vdsm service for running helper functions as root
Requires=libvirtd.service
Wants=lldpad.service
After=lldpad.service libvirtd.service

[Service]
Type=simple
LimitCORE=infinity
EnvironmentFile=-/etc/sysconfig/supervdsmd
ExecStart=@LIBEXECDIR@/daemonAdapter "@LIBEXECDIR@/supervdsmd" --sockfile "@VDSMRUNDIR@/svdsm.sock"
Restart=always
