# static/usr/lib/systemd/system/vdsm-network.service.in
# SPDX-FileCopyrightText: Red Hat, Inc.
# SPDX-License-Identifier: GPL-2.0-or-later

[Unit]
Description=Virtual Desktop Server Manager network restoration
Wants=network.target
Requires=libvirtd.service openvswitch.service NetworkManager.service
After=libvirtd.service openvswitch.service NetworkManager.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/sysconfig/vdsm
ExecStartPre=@BINDIR@/vdsm-tool dump-bonding-options
ExecStart=@BINDIR@/vdsm-tool restore-nets
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
