# static/usr/lib/systemd/system/vdsmd.service.in
[Unit]
Description=Virtual Desktop Server Manager
Requires=multipathd.service time-sync.target \
         iscsid.service rpcbind.service supervdsmd.service sanlock.service \
         vdsm-network.service
After=multipathd.service iscsid.service rpcbind.service supervdsmd.service \
      sanlock.service vdsm-network.service
Before=libvirt-guests.service
Wants=mom-vdsm.service ovirt-imageio-daemon.service abrtd.service \
      dev-hugepages1G.mount libvirt-guests.service kdump.service

[Service]
Type=simple
LimitCORE=infinity
EnvironmentFile=-/etc/sysconfig/vdsm
ExecStartPre=@LIBEXECDIR@/vdsmd_init_common.sh --pre-start
ExecStart=@VDSMDIR@/daemonAdapter -0 /dev/null -1 /dev/null -2 /dev/null "@VDSMDIR@/vdsmd"
# WARNING: Do not use ExecStopPost, breaks killing child processes.
# https://bugzilla.redhat.com/1761260
Restart=always
Nice=-20
User=@VDSMUSER@
Group=@VDSMGROUP@
PermissionsStartOnly=true
TimeoutStopSec=@SERVICE_STOP_TIMEOUT@
KillMode=mixed
LimitNOFILE=4096
@SYSTEMD_LIMIT_MEMLOCK@

[Install]
WantedBy=multi-user.target
