# static/usr/lib/systemd/system/vdsm-network-init.service.in
[Unit]
Description=Virtual Desktop Server Manager network IP+link restoration
Wants=network.target
After=@OPENVSWITCHSERVICE@
Before=libvirtd.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/sysconfig/vdsm
ExecStart=@BINDIR@/vdsm-tool dump-bonding-options
ExecStart=@BINDIR@/vdsm-tool restore-nets-init
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
