# vdsm/supervdsmd.service.in
[Unit]
Description="Auxiliary vdsm service for running helper functions as root"
Requires=libvirtd.service

[Service]
Type=simple
ExecStart=@VDSMDIR@/daemonAdapter "@PYTHON@" "@VDSMDIR@/supervdsmServer.pyc" --sockfile "@VDSMRUNDIR@/svdsm.sock"
Restart=on-abort
