# init/systemd/supervdsmd.service.in
[Unit]
Description="Auxiliary vdsm service for running helper functions as root"
Requires=libvirtd.service

[Service]
Type=simple
EnvironmentFile=-/etc/sysconfig/supervdsmd
ExecStart=@VDSMDIR@/daemonAdapter "@VDSMDIR@/supervdsmServer" --sockfile "@VDSMRUNDIR@/svdsm.sock"
Restart=on-abort
