# init/systemd/vdsm-network.service.in
[Unit]
Description=Virtual Desktop Server Manager network restoration
Wants=network.target
Requires=libvirtd.service
After=libvirtd.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/sysconfig/vdsm
ExecStartPre=@BINDIR@/vdsm-tool --vvverbose --append --logfile=@VDSMLOGDIR@/upgrade.log upgrade-unified-persistence
ExecStartPre=@BINDIR@/vdsm-tool --vvverbose --append --logfile=@VDSMLOGDIR@/upgrade.log upgrade-3.0.0-networks
ExecStart=@BINDIR@/vdsm-tool restore-nets
KillMode=process
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
