# init/systemd/vdsm-network-init.service.in
[Unit]
Description=Virtual Desktop Server Manager network IP+link restoration
Wants=network.target
# do not require openvswitch until it starts on AMD (BZ#1378501)
Wants=@OPENVSWITCHSERVICE@
After=@OPENVSWITCHSERVICE@
Before=libvirtd.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/sysconfig/vdsm
ExecStart=@BINDIR@/vdsm-tool restore-nets-init
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
