# vdsm_hooks/localdisk/12-vdsm-localdisk.rules
# SPDX-FileCopyrightText: Red Hat, Inc. and/or its affiliates.
# SPDX-License-Identifier: GPL-2.0-or-later

# Vdsm udev rules for local lvs

# "add" event is processed on coldplug only, so we need "change", too.
ACTION!="add|change", GOTO="lvm_end"

# Filter out vgs which do not look like an ovirt local vg
ENV{DM_VG_NAME}!="ovirt-local", GOTO="lvm_end"

# Volumes used as local disks lvs
# WARNING: we cannot use OWNER, GROUP and MODE since using any of them will
# change the selinux label to the default, causing vms to pause after refreshing the lv.
# See https://bugzilla.redhat.com/1147910
ENV{DM_LV_NAME}=="*", RUN+="/usr/bin/chown vdsm:qemu $env{DEVNAME}", GOTO="lvm_end"

LABEL="lvm_end"
