# init/upstart/supervdsmd.upstart.in
# supervdsmd - Auxiliary vdsm service for running helper functions as root
#

description "Auxiliary vdsm service for running helper functions as root"

start on runlevel [2345] and started libvirt-bin
stop on runlevel [!2345]
chdir "@VDSMDIR@"
console log
respawn

# All commands called inside this script section except the daemon itself
# should not fork, otherwise Upstart traces the wrong pid.
# ".", "[", "&&" are built-in command or key-word, no fork.
# bash exec does not fork, just execve the target binary.
# So no "expect" stanza is needed.
script
    [ -f "/etc/default/supervdsmd" ] && . "/etc/default/supervdsmd"
    exec "@VDSMDIR@/daemonAdapter" "@VDSMDIR@/supervdsmServer" --sockfile "@VDSMRUNDIR@/svdsm.sock"
end script
